# cron/abuseblacklist.crontab
# AbuseBlacklist daily update cron
# Runs every day at 03:17 UTC/server time.
#
# Install example:
#   crontab cron/abuseblacklist.crontab
#
# If the server timezone is not UTC, adjust the hour/minute or set CRON_TZ if supported.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

17 3 * * * cd /opt/AbuseBlacklist && /usr/bin/env bash scripts/run_daily_update.sh >> /var/log/abuseblacklist-daily.log 2>&1
